# snapserver-addon/rootfs/etc/services.d/snapserver/run
#!/command/with-contenv bashio
# ==============================================================================
# Runs Snapserver
# ==============================================================================
bashio::log.info "Starting Snapserver..."

# Run Snapserver
exec snapserver -c /config/snapserver.conf
